# buildroot/profiles/director-aussie/target_skeleton/etc/inittab
# /etc/inittab
#
# Copyright (C) 2001 Erik Andersen <andersen@codepoet.org>
#
# Note: BusyBox init doesn't support runlevels.  The runlevels field is
# completely ignored by BusyBox init. If you want runlevels, use
# sysvinit.
#
# Format for each entry: <id>:<runlevels>:<action>:<process>
#
# id        == tty to run on, or empty for /dev/console
# runlevels == ignored
# action    == one of sysinit, respawn, askfirst, wait, and once
# process   == program to run

# Startup the system
# now run any rc scripts
::sysinit:/etc/init.d/rcS
/dev/ttyS0::shutdown:/etc/init.d/rcStop
/dev/ttyS0::respawn:/bin/cttyhack /bin/login.sh
/dev/ttyS0::respawn:/bin/tacmon.sh
/dev/ttyS0::respawn:/bin/tsyslogd -r -h -n --rotate=7
/dev/ttyS0::respawn:/bin/rhttpd
/dev/ttyS0::respawn:/bin/zapd
/dev/ttyS0::respawn:/bin/clusterD



